# .config/systemd/user/batsignal.service.d/options.conf
[Service]
ExecStart=
ExecStart=batsignal -m 30 -p -d 5 -c 10 -w 30 -I /home/toms/.config/system_resources/icons8-android-l-battery-64.png
